# internal/filesystem/etc/systemd/system/charge-limit-reapply.service
[Service]
Type=oneshot
ExecStart=/usr/local/bin/charge-limit-reapply.sh
