# internal/filesystem/etc/systemd/system/charge-limit-reapply.timer
[Timer]
OnBootSec=1s
OnUnitActiveSec=8s
AccuracySec=1s

[Install]
WantedBy=timers.target
