# internal/filesystem/home/user/.config/systemd/user/hyprland-session.target
[Unit]
Description=Hyprland session wrapper
BindsTo=graphical-session.target
Wants=graphical-session.target
Before=graphical-session.target
