# internal/filesystem/home/user/.config/systemd/user/usb-remove.service
[Unit]
Description=Play sound on USB removal

[Service]
Type=oneshot
ExecStart=/usr/bin/pw-play %h/.local/share/sounds/usb_remove.ogg
